# .config/autostart/startup.sh.desktop
[Desktop Entry]
Exec=/home/kaya-sem/dotfiles/scripts/startup.sh
Icon=dialog-scripts
Name=startup.sh
Type=Application
X-KDE-AutostartScript=true
